# gpio-server.service
[Unit]
Description=GPIO Server for Telephone

[Service]
ExecStart=/home/hans/fon/gpio-server

[Install]
WantedBy=multi-user.target
